# systemd/adhan.service
[Unit]
Description=Adhan appliance
# time-sync.target: a Raspberry Pi has no RTC, so wait until the clock is
# NTP-synced before scheduling — otherwise prayers scheduled against the wrong
# boot-time clock never fire (see install.sh, which enables systemd-time-wait-sync).
After=network-online.target time-sync.target adhan-bt-keepalive.service
Wants=network-online.target time-sync.target

[Service]
Type=simple
User=adhan
Group=adhan
SupplementaryGroups=audio bluetooth
ExecStart=/opt/adhan/.venv/bin/adhan --config /etc/adhan/config.yaml --media /etc/adhan/media --state /var/lib/adhan/state.json run
Restart=always
RestartSec=10
StateDirectory=adhan

[Install]
WantedBy=multi-user.target
